# cluster/base@.container
[Container]
Image="localhost/%j:latest"
HostName="%j-%i"
