# cluster/base@.container
[Container]
Image="localhost/%j:latest"
HostName="%j-%i"
Network=lab
#Network=cluster.network
NetworkAlias=%j
