# cluster/base@.container
[Container]
Image="localhost/%j:latest"
Network=host
HostName="%j-%i"
