# services/archive-logs.timer
[Unit]
Description=archive logs to remote timer

[Timer]
OnCalendar=daily
Unit=archive-logs.service

[Install]
WantedBy=timers.target
