# services/base-.service
[Unit]
OnFailure=notify-slack@%i.service
StartLimitIntervalSec=30
StartLimitBurst=2

[Service]
LimitCORE=infinity
Restart=on-failure
