# services/base-job@.service
[Container]
Image=ghcr.io/avagordon01/prod:latest
Network=host
AutoUpdate=registry

[Service]
MemoryHigh=2G
MemoryMax=4G
