# services/base-strategy-shanghai@4.timer
[Timer]
OnCalendar=mon..fri 08:30 Asia/Shanghai
