# services/podman-clean.service
[Unit]
Description=podman clean service

[Service]
Type=oneshot
ExecStart=/etc/systemd/systemd/podman-clean.sh
