# cluster-systemd/cluster.slice
[Unit]
Wants=\
    control@0.service \
    farm@0.service \
    farm@1.service \
    farm@2.service \
    farm@3.service \
    prod@newyork.service \
    prod@chicago.service \
    prod@taiwan.service \
    prod@shanghai.service
